# ION/09_integrations/systemd/user/ion-chatops.service.template
[Unit]
Description=ION ChatOps Local Bridge
After=network.target

[Service]
Type=simple
WorkingDirectory=__ION_ROOT__
Environment=PYTHONDONTWRITEBYTECODE=1
Environment=PYTHONPATH=__ION_PYTHONPATH__
ExecStart=/usr/bin/python3 -S -m kernel.ion_chatops_bridge --ion-root . --host 127.0.0.1 --port 8767 --serve
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=default.target
